# roles/docker-compose/files/systemd/user/rebuild.service
[Unit]
Description=Docker Compose Rebuild

[Service]
Type=oneshot
ExecStart=/home/vagrant/bin/rebuild.sh

