# roles/syncing/files/systemd/user/dockerdevsync.service
[Unit]
Description=Docker Dev sync

[Service]
Type=oneshot
ExecStart=/home/vagrant/bin/dockersync.sh

