# localsystemd/dockerdevsync.timer
[Unit]
Description=Docker Dev Sync Timer

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m
Unit=dockerdevsync.service

[Install]
WantedBy=timers.target
