# localsystemd/dockerdevsync.timer
[Unit]
Description=Docker Dev Sync Timer

[Timer]
OnBootSec=60m
OnUnitActiveSec=60m
Unit=dockerdevsync.service

[Install]
WantedBy=timers.target
