# localsystemd/dockerdevsync.service
[Unit]
Description=Docker Dev sync

[Service]
Type=oneshot
ExecStart=/home/mark/bin/dockerdevsync.sh

